# utils/openRC
#!/sbin/openrc-run
# Description: SquidStats Web Application

name="SquidStats"
description="SquidStats Web Application"
command="/opt/SquidStats/venv/bin/python3"
command_args="/opt/SquidStats/app.py"
directory="/opt/SquidStats"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
    need net
    after firewall
}

start_pre() {
    if [ -f "/opt/SquidStats/.env" ]; then
        set -a
        . "/opt/SquidStats/.env"
        set +a
    fi
    export PATH="/opt/SquidStats/venv/bin:$PATH"
}